# systemd/newsroom-pipeline.service
[Unit]
Description=AI Newsroom Pipeline (one-shot run)
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/ai-newsroom/agent
EnvironmentFile=%h/ai-newsroom/agent/.env
ExecStart=/bin/bash %h/ai-newsroom/scripts/run-pipeline.sh
StandardOutput=journal
StandardError=journal
